# CTF_1/ctf/oops-ctf/oops-ctf/oops1/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:mariadb]
command=/usr/sbin/mariadbd --user=mysql
autorestart=true

[program:php-fpm]
command=/bin/bash -c "mkdir -p /run/php && exec $(ls /usr/sbin/php-fpm* | head -1) -F"
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true
